# services/xclaw-agentguard.service
[Unit]
Description=XClaw AgentGuard Protection Engine
After=network.target

[Service]
Type=simple
ExecStart=xclaw-agentguard engine-start --daemon
Restart=always

[Install]
WantedBy=multi-user.target